# board/picopi/overlay/etc/udev/rules.d/80-fbcp.rules
SUBSYSTEM=="graphics" ACTION=="add" ENV{DEVNAME}=="/dev/fb1", RUN+="/sbin/start-stop-daemon start-stop-daemon -S -b -x /usr/sbin/fbcp"
