# board/picopi/overlay/etc/udev/rules.d/90-network.rules
SUBSYSTEM=="net", ACTION=="add", ENV{INTERFACE}!="lo", RUN+="/sbin/udhcpc -b -t 1 -i %k"
